# references/launchd/ai.openclaw.semiconductor-afterclose.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.semiconductor-morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/your-user/.openclaw/scripts/semiconductor-daily-update.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SEMI_ENFORCE_SCHEDULE_WINDOW</key>
        <string>1</string>
        <key>SEMI_SCHEDULE_MODE</key>
        <string>after_us_close_full</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>4</integer>
            <key>Minute</key>
            <integer>12</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>12</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>45</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/your-user/.openclaw/logs/semiconductor-daily-update.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/your-user/.openclaw/logs/semiconductor-daily-update.log</string>
</dict>
</plist>
